# systemd/xboard-report.service
[Unit]
Description=Report official Xray-core user traffic to XBoard
After=network-online.target docker.service xboard-sync.service
Wants=network-online.target docker.service

[Service]
Type=simple
WorkingDirectory=/opt/xray-sync
ExecStart=/bin/bash -c 'while true; do /usr/bin/python3 /opt/xray-sync/xboard_report.py; sleep 60; done'
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
