# systemd/xboard-sync.service
[Unit]
Description=XBoard to official Xray-core config sync
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
WorkingDirectory=/opt/xray-sync
ExecStart=/usr/bin/python3 /opt/xray-sync/xboard_sync.py
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
